# docker/repeka/crontab
* * * * * root su - www-data -s /bin/bash -c "/usr/local/bin/php /var/www/html/bin/console repeka:cyclic-tasks:dispatch" >/var/log/cron.out 2>&1
